# lib/systemd/system/lastore-daemon.service
[Unit]
Description=Deepin Lastore Daemon
After=display-manager.service

[Service]
Type=dbus
BusName=org.deepin.dde.Lastore1
ExecStart=/usr/libexec/lastore-daemon/lastore-daemon
StandardOutput=null
StandardError=null
StateDirectory=lastore
CacheDirectory=lastore
RuntimeDirectory=lastore
RuntimeDirectoryMode=0750