# lib/systemd/system/lastore-smartmirror-daemon.service
[Unit]
Description=Deepin Lastore SmartMirror Daemon
After=display-manager.service

[Service]
Type=dbus
BusName=com.deepin.lastore.Smartmirror
ExecStart=/usr/bin/lastore-smartmirror-daemon