# lib/systemd/system/lastore-daemon.service
[Unit]
Description=Deepin Lastore Daemon
After=display-manager.service


[Service]
Type=dbus
BusName=com.deepin.lastore
ExecStart=/usr/libexec/lastore-daemon/lastore-daemon
StandardOutput=null
StandardError=null
Environment=GOTRACEBACK=crash