# lib/systemd/system/lastore-update-metadata-info.timer
[Unit]
Description=Update lastore metadata info

[Timer]
OnUnitActiveSec=3h
OnBootSec=15min

[Install]
WantedBy=timers.target
