# lib/systemd/system/lastore-daemon.service
[Unit]
Description=Deepin Lastore Daemon
After=display-manager.service

[Service]
Type=dbus
BusName=org.deepin.dde.Lastore1
ExecStart=/usr/libexec/lastore-daemon/lastore-daemon

Restart=on-failure
RestartForceExitStatus=2
StartLimitBurst=0

StandardOutput=null
StandardError=journal

StateDirectory=lastore
CacheDirectory=lastore