# etc/xdg/autostart/lastore-session-helper.desktop
[Desktop Entry]
Type=Application
Name=lastore-session-helper
Exec=/usr/bin/lastore-session-helper
TryExec=/usr/bin/lastore-daemon
NoDisplay=true

